# crontab
0 11,18 * * * cd /mit/locker/cron_scripts && bash police-logs.sh >/dev/null 2>&1
